# disk-monitor/disk-collector.service
[Unit]
Description=Disk Health Data Collector
After=local-fs.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=jeb
WorkingDirectory=/home/jeb/programs/gemini_cli_workspace/disk-monitor
# Use the self-contained virtualenv
ExecStart=/home/jeb/programs/gemini_cli_workspace/disk-monitor/venv/bin/python3 collect_disk_features.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
